# yabai/Library/LaunchAgents/com.eduardo.yabai-bootstrap.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.eduardo.yabai-bootstrap</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>exec "$HOME/.config/yabai/login-bootstrap.sh"</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>ProcessType</key>
  <string>Interactive</string>

  <key>StandardOutPath</key>
  <string>/tmp/yabai-login-bootstrap.out.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/yabai-login-bootstrap.err.log</string>
</dict>
</plist>
